# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; ── PHP-FPM ──────────────────────────────────────────────────
[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Laravel Queue Worker ──────────────────────────────────────
[program:laravel-queue]
command=php /var/www/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --queue=default
directory=/var/www
user=www-data
autostart=true
autorestart=true
numprocs=2
process_name=%(program_name)s_%(process_num)02d
priority=10
stopwaitsecs=600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Laravel Scheduler ────────────────────────────────────────
[program:laravel-scheduler]
command=/bin/sh -c "while true; do php /var/www/artisan schedule:run --no-interaction; sleep 60; done"
directory=/var/www
user=www-data
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
